# ansible/roles/app-config/templates/csv-app.service.j2
[Unit]
Description={{ app_name }} Web Application
After=network.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/opt/csv-app
EnvironmentFile=/opt/csv-app/config/app.env
ExecStart=/usr/local/bin/gunicorn --config /opt/csv-app/config/gunicorn.conf.py app:app
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
